# docker/crontab
0 */12 * * * /usr/local/bin/php /var/www/html/bin/console app:zotero:load >> /var/log/cron.log 2>&1
